# logrotate.cron
*/5 *  *  *  * logrotate /app/logrotate.conf
